# fleet/services/fleet-health-monitor.service
[Unit]
Description=Fleet Health Monitor — 24/7 service watcher, self-heals, alerts Casey
After=network.target
Wants=plato.service

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/home/ubuntu/.openclaw/workspace/repos/fleet-health-monitor
ExecStart=/usr/bin/npx ts-node src/index.ts
Restart=always
RestartSec=10
StandardOutput=append:/tmp/fleet-health-monitor.log
StandardError=append:/tmp/fleet-health-monitor.log
Environment=PLATO_SERVER=http://localhost:8847

[Install]
WantedBy=multi-user.target
