# fleet/services/intent-inference.service
[Unit]
Description=Intent Inference Engine — reverse-actualization intent classifier
After=network.target
Wants=plato.service

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/home/ubuntu/.openclaw/workspace/repos/intent-inference
ExecStart=/usr/bin/npx ts-node src/index.ts
Restart=always
RestartSec=10
StandardOutput=append:/tmp/intent-inference.log
StandardError=append:/tmp/intent-inference.log
Environment=PLATO_SERVER=http://localhost:8847

[Install]
WantedBy=multi-user.target
